# dist/macos/styx-receiver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ghreprimand.styx-receiver</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Styx Receiver.app/Contents/MacOS/styx-receiver</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>3</integer>
    <key>AssociatedBundleIdentifiers</key>
    <array>
        <string>com.ghreprimand.styx-receiver</string>
    </array>
    <key>StandardOutPath</key>
    <string>/tmp/styx-receiver.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/styx-receiver.stderr.log</string>
</dict>
</plist>
